# deploy/systemd/kitchen-bot.service
[Unit]
Description=Register Kitchen MAX bot webhook
After=network.target kitchen-api.service

[Service]
Type=oneshot
WorkingDirectory=/opt/kitchen-platform
EnvironmentFile=/opt/kitchen-platform/.env
ExecStart=/opt/kitchen-platform/.venv/bin/python -m bot.main --register-webhook
RemainAfterExit=yes
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
