# deploy/systemd/kitchen-broadcast-worker.service
[Unit]
Description=Kitchen platform broadcast worker
After=network.target kitchen-api.service

[Service]
Type=simple
WorkingDirectory=/opt/kitchen-platform
EnvironmentFile=/opt/kitchen-platform/.env
ExecStart=/opt/kitchen-platform/.venv/bin/python worker/broadcast_worker.py
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
